# Docker-Choose/nginx/supervisor/memcache.conf
[program:memcache]
command=memcached -u %(ENV_MEMCACHED_USER)s -vp 11211
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
